# bin/deploy/crontab
0  0  *   *   *   rails     rails runner RefreshStaticRequestsJob.perform_now
