# bin/deploy/crontab
# vim: filetype=crontab
# m h dom mon dow  user      command
0  *   *   *   *   rails     $HOME/bin/cron rails runner RefreshStaticRequestsJob.perform_now
